# bambuddy/rootfs/etc/services.d/bambuddy/finish
#!/usr/bin/with-contenv bashio

echo "BambuBuddy stopped"