# bambuddy/rootfs/etc/services.d/bambuddy/run
#!/usr/bin/with-contenv bashio

echo "Starting BambuBuddy..."

cd /app
exec python3 main.py