# bambuddy/rootfs/etc/services.d/bambuddy/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH=/data
PORT=$(bashio::config 'port')
LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Starting BambuBuddy on port ${PORT}"

cd /app

exec python3 main.py --port ${PORT}